# crontab

13 4 * * *  www-data	php /var/www/omkt/modules/api/scripts/cron.daily.php 



