# crontab



13 4 * * *   /var/www/omkt/modules/api/scripts/cron.daily.php

